# deb/moidemd.service.in
[Unit]
Description=moidemd elevation service
Requires=network.target
After=multi-user.target

[Service]
Type=simple
Restart=always
RestartSec=30
StartLimitInterval=0
User=moidemd
ExecStart=/usr/sbin/moidemd /var/lib/moidemd/%%DEM%%

[Install]
WantedBy=multi-user.target
